# osbuild-manifests/files/etc/systemd/system/neptune-accounts-service-init.service
[Unit]
Description=Init neptune accounts-service conf
ConditionPathExists=!/var/lib/AccountsService/users/neptune
Before=gdm.service

[Service]
Type=service
ExecStart=/bin/sh -c ' \
  echo -e "[User]\nSession=neptune-wayland\nXSession=neptune-xorg\nPasswordHint=\nIcon=/var/lib/AccountsService/icons/guest\nSystemAccount=false\n" \
  > /var/lib/AccountsService/users/neptune'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target